# debian/scrapyd.upstart
description "scrapyd - run Scrapy spiders"

start on runlevel [2345]
stop on runlevel [016] or unmounting-filesystem

chdir /var/lib/scrapyd

script
  [ -r /etc/default/scrapyd ] && . /etc/default/scrapyd
  exec /usr/bin/scrapyd -u scrapy -g nogroup -l /var/log/scrapyd/scrapyd.log
end script
